# roles/etcd/files/etcd.service.d/20-cloudinit.conf
[Service]
Environment="ETCD_ADDR=172.17.8.101:4001"
Environment="ETCD_DISCOVERY=https://discovery.etcd.io/75d59bd75fcd8dac12102d7ab9f94753"
Environment="ETCD_NAME=d7ab8ee6a078471da81e38cb96dcd357"
Environment="ETCD_PEER_ADDR=172.17.8.101:7001"
